# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:lightpanda]
command=/usr/local/bin/lightpanda serve --host 0.0.0.0 --port 9222 --log_level info
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:cdp-proxy]
command=/usr/local/bin/cdp-proxy
environment=PORT="9333",BROWSER_ENDPOINTS="ws://127.0.0.1:9222",BROWSER_COUNT="1",MAX_PER_BROWSER="30"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
startsecs=3

[program:linkedin-api]
command=python linkedin_finder.py
directory=/app
environment=CDP_PROXY_URL="http://127.0.0.1:9333",WORKERS="1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
startsecs=5
